# Dockerfiles/rtorrent/template/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
priority=1
command=/root/startup-nginx.sh
autorestart=false

[program:php5-fpm]
priority=3
command=/root/startup-php.sh

